# museu
#!/sbin/openrc-run

directory=/var/www/$RC_SVCNAME/
pidfile=/run/$RC_SVCNAME/$RC_SVCNAME.pid

depend(){
    need net
}

start() {
    ebegin "Starting ${RC_SVCNAME}"
    start-stop-daemon --start \
        -d "${directory}" \
        --exec /usr/bin/node src/index.js \
        --pidfile ${pidfile} \
        --background
    eend $?
}

stop() {
    ebegin "Stopping ${RC_SVCNAME}"
    start-stop-daemon --stop \
        --pidfile ${pidfile} \
    eend $?
}

reload() {
    ebegin "Restarting ${RC_SVCNAME}"
        kill -HUP ${pidfile}
        start
    eend $?
}